# waycast-daemon.service
[Unit]
Description=Waycast application launcher daemon
Documentation=https://waycast.dev
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=%i/bin/waycast-daemon
Restart=on-failure
RestartSec=5
# Ensure the daemon starts in the user's environment
Environment=XDG_RUNTIME_DIR=%t
Environment=WAYLAND_DISPLAY=wayland-0

[Install]
WantedBy=default.target